# systemd/user/i3@blocks_socket@pulseaudio.service.d/override.conf
[Unit]
After=pulseaudio.service
Requisite=pulseaudio.service
