# profile.d/android-studio.sh
# Android studio
export ANDROID_SDK_HOME="${XDG_CONFIG_HOME}/android"
