# profile.d/arch-repo.sh
# Archrepo configuration
export SBP_REPO_NAME="sbp"
export SBP_REPO_URL='s3://sbp-arch/repo'
export SBP_REPO_SOURCE="${XDG_DATA_HOME}/arch"
