# profile.d/base16.sh
# Status bar colors: base16-default-dark
export base00='#181818'
export base01='#282828'
export base02='#383838'
export base03='#585858'
export base04='#b8b8b8'
export base05='#d8d8d8'
export base06='#e8e8e8'
export base07='#f8f8f8'
export base08='#ab4642'
export base09='#dc9656'
export base0A='#f7ca88'
export base0B='#a1b56c'
export base0C='#86c1b9'
export base0D='#7cafc2'
export base0E='#ba8baf'
export base0F='#a16946'
