# profile.d/bspwm.sh
#!/bin/sh

# BSPWM variable names
export SBP_XPAPER_DIR="${HOME}/Pictures/Wallpapers"
