# profile.d/defaults.sh
export EDITOR="/usr/bin/nvim"
export BROWSER="/usr/bin/firefox"
export TERMINAL="/usr/bin/alacritty"
