# profile.d/defaults.sh
export EDITOR="/usr/bin/nvim"
export BROWSER="/usr/bin/qutebrowser"
export TERMINAL="/usr/bin/kitty"
