# profile.d/hardware.sh
# Monitors
if [ "$(uname -n)" = 'sbp-arch-laptop' ] ; then
    export BRI_SCR="intel_backlight"
    export BRI_KBD="system76_acpi::kbd_backlight"
fi
