# profile.d/mpd.sh
# MPD stuff
export MPD_HOST='localhost'
