# profile.d/neomutt.sh
# Notmuch
export NOTMUCH_CONFIG="${XDG_CONFIG_HOME}/neomutt/notmuch-config"
