# profile.d/openfoam.sh
# OpenFOAM
export FOAM_INST_DIR="/opt/OpenFOAM"
