# profile.d/pass.sh
# Password store location
export PASSWORD_STORE_DIR="${HOME}/.local/pass"
