# profile.d/pass.sh
# Password store location
export PASSWORD_STORE_DIR="${HOME}/.pass"
export PASSWORD_STORE_GIT="${HOME}/.pass"
