# profile.d/rofi-pass.sh
# Rofi-pass config
export ROFI_PASS_CONFIG="${XDG_CONFIG_HOME}/rofi/pass.conf"
