# profile.d/sbp.sh
# Archrepo configuration
export SBP_REPO_NAME="sbp"
export SBP_REPO_URL='s3://sbp-arch/repo'
export SBP_REPO_DIR="${HOME}/Projects/Archlinux"

# Userspace script variables
export SBP_WPAPER_DIR="Wallpapers"
