# profile.d/sbpman.sh
# Define Zim location
export SBPMAN_HOME="${HOME}/Projects/Arch/"
