# profile.d/steam.sh
# Just steam things
export LD_PRELOAD="${LD_PRELOAD}:/usr/lib/libgio-2.0.so.0:/usr/lib/libglib-2.0.so.0"
export STEAM_RUNTIME=0
