# profile.d/sysinfo.sh
#!/bin/sh
# Define these to change behaviour of sysinfo scripts
export SYSINFO_CAL_LOC="${HOME}/Documents/Calendar"
export SYSINFO_CAL_MAXLEN=15
export SYSINFO_CPU_POLL=5
export SYSINFO_DATE_STYLE='reg'
export SYSINFO_DISK_POLL=900
export SYSINFO_EMAIL_LOC="${HOME}/Documents/Mail/Gmail"
export SYSINFO_MEM_POLL=30
export SYSINFO_MPD_POLL=30
export SYSINFO_MPD_MAXLEN=20
export SYSINFO_TODO_LOC="${HOME}/Documents/Todo"
export SYSINFO_TODO_MAXLEN=15
