# profile.d/wayland.sh
# Wayland related config files
export SBP_WPAPER_DIR="${HOME}/Pictures/Wallpapers"
