# profile.d/x11.sh
# Location of X11 config files
export XINITRC="${XDG_CONFIG_HOME}/X11/clientrc"
export XSERVERRC="${XDG_CONFIG_HOME}/X11/serverrc"
# User variable names
export SBP_XPAPER_DIR="${HOME}/Pictures/Wallpapers"
