# profile.d/zimfw.sh
# Define Zim location
export ZIM_HOME="${XDG_DATA_HOME}/zim"
