# autostart/autorandr.desktop
[Desktop Entry]
Exec=/usr/bin/autorandr --change
OnlyShowIn=bspwm
