# autostart/polybar.desktop
[Desktop Entry]
Name=Polybar
Comment=Statusbar for BSPWM
Exec=polybar bspwm
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=bspwm;
