# autostart/rambox.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=rambox
Comment=ramboxstartup script
Exec=/opt/rambox/rambox
StartupNotify=false
Terminal=false