# autostart/sxhkd.desktop
[Desktop Entry]
Name=SXHKD
Comment=Simple X HotKey Daemon for BSPWM
Exec=sxhkd
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=bspwm;
