# systemd/user/i3@blocks_socket@mpd.service.d/override.conf
[Unit]
After=mpd.service
Requisite=mpd.service
