# arch/bbaserdem-devel.install
# Maintained by bbaserdem-base
post_install() {
    post_upgrade()
}

post_upgrade() {
    systemctl enable --now libvirtd.service
}
