# systemd/user/blueman-applet@.service
[Unit]
Description=Bluetooth manager applet on X display %I
Documentation=man:blueman-applet(1)

[Service]
Environment=DISPLAY=%I
Type=simple
ExecStart=/usr/bin/blueman-applet
SyslogIdentifier=blueman-applet

[Install]
WantedBy=systray@.service
