# systemd/user/feh.timer
[Unit]
Description=Random wallpaper change on interval

[Timer]
OnUnitActiveSec=15m
Unit=feh.service

[Install]
WantedBy=WM.target
