# systemd/user/systray@.target
[Unit]
Description=Minimized gui apps to X11 display %I
After=wm-running@%i.target
