# systemd/user/systray@.target
[Unit]
Description=Minimized gui apps to X11 display %I
