# systemd/user/wm-running@.target
[Unit]
Description=Window Manager is running on X display %I
