# systemd/user/WM.target
[Unit]
Description=Window manager is running
