# systemd/user/mpdscribble.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/mpdscribble --no-daemon --conf %h/.cache/mpdscribble/mpdscribble.conf
