# packagebuilds/arch/devel.install
# Maintained by sbp-devel
post_install() {
    post_upgrade()
}

post_upgrade() {
    systemctl enable --now libvirtd.service
}
