# packagebuilds/arch/laptop.install
# Maintained by bbaserdem-laptop
post_install() {
    mkdir /opt/Steam
    # TLP service
    systemctl enable --now tlp.service
    systemctl enable --now tlp-sleep.service
    # TLP and NVME fix 
    sed -i -e '/SATA_LINKPWR_ON_BAT/s/=".*"/="max_performance"/' /etc/default/tlp
    post_upgrade()
}

post_upgrade() {
}
