# packagebuilds/arch/workstation.install
# Maintained by sbp-workstation
post_install() {

    # Enable SSH
    echo "Enabling SSH socket . . ."
    systemctl enable --now sshd.socket

    post_upgrade()
}

post_upgrade() {
}
